# personal/aurora.service
[Unit]
Description=Aurora daemon service

[Service]
Type=oneshot
ExecStart=/usr/bin/python /home/yyannick/Aurora/daemon.py